# cb4/cb4-lubuntu-desktop-hdmi/overlay/etc/udev/rules.d/11-block-detec.rules
KERNEL=="sd[a-z][0-9]",SUBSYSTEM=="block",ACTION=="add", RUN+="/bin/cubie-mount.sh %k"
KERNEL=="mmcblk1p[0-9]",ACTION=="add", RUN+="/bin/cubie-mount.sh %k"
